# deploy/fiam-dashboard.service
[Unit]
Description=fiam dashboard JSON API backend
After=network.target

[Service]
Type=simple
User=fiet
Group=fiet
WorkingDirectory=/home/fiet/fiam-code
EnvironmentFile=-/home/fiet/fiam-code/.env
Environment=FIAM_CC_TRANSPORT=channel
ExecStart=/home/fiet/fiam-code/.venv/bin/python /home/fiet/fiam-code/scripts/dashboard_server.py --port 8766
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Harden: read-only for most paths, read-write only for logs
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/fiet/live /home/fiet/fiam-code/logs /home/fiet/fiam-code/store /home/fiet/.claude /home/fiet/.claude.json /home/fiet/.cache /home/fiet/.config /tmp
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
